# visitorkiosk-fe.service
[Unit]
Description=Visitor Kiosk Frontend (Vite Preview)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/VisitorKiosk/fe
Environment=NODE_ENV=production

# Vite preview serves the built /dist output
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 8007
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
